# etc/cron.d/crontab.root
0 */6 * * * touch /root/.weather_update >/dev/null 2>&1
