# etc/update-motd.d/20-weather
#!/bin/sh

[ ! -f /root/.weather_update ] && exit 0
weather
rm /root/.weather_update
