# etc/update-motd.d/21-rainbow-fortune-cow
#!/bin/sh

/usr/games/fortune | /usr/games/cowsay -T U | lolcat
